# root/etc/s6-overlay/s6-rc.d/rabbitmq/run
#!/command/with-contenv sh

export HOME=/var/lib/rabbitmq
export PATH=$PATH:/opt/rabbitmq/sbin

echo "Starting rabbitmq-server..."
s6-setuidgid unms rabbitmq-server
