# root/etc/services.d/siridb/run
#!/command/with-contenv sh

echo "Starting siridb-server..."
s6-setuidgid unms siridb-server
