# root/etc/services.d/rabbitmq/run
#!/usr/bin/with-contenv sh

export HOME=/var/lib/rabbitmq
export PATH=$PATH:/opt/rabbitmq/sbin

echo "Starting rabbitmq-server..."
exec /opt/rabbitmq/sbin/rabbitmq-server
